# src/crontab
# m h  dom mon dow   command
03  04 *   *   *     ( cd $HOME; date; ./make_all.sh ; date )
